# deploy/systemd/cnb-watcher.timer
[Unit]
Description=Poll CNB approved production release every minute

[Timer]
OnBootSec=30s
OnUnitInactiveSec=60s
RandomizedDelaySec=10s
Persistent=true
Unit=cnb-watcher.service

[Install]
WantedBy=timers.target
